# assets/etc/systemd/system/r2.timer
# Created by Alex Thomson

[Unit]
Description=Run r2 on boot

[Timer]
OnBootSec=30
Unit=r2.service

[Install]
WantedBy=timers.target